# backups/version_2.4.0_beta_20250901_193155/templates/service-orchestration-monitor.service.j2
[Unit]
Description=Service Orchestration Monitoring
After=network.target service-orchestrator.service
Requires=service-orchestrator.service

[Service]
Type=simple
ExecStart=/usr/local/bin/service_orchestration_monitor.sh
StandardOutput=append:/var/log/service_orchestration/monitor.log
StandardError=append:/var/log/service_orchestration/monitor.log

[Install]
WantedBy=multi-user.target 