# backups/version_2.4.0_beta_20250901_193155/templates/service-orchestration-monitor.timer.j2
[Unit]
Description=Service Orchestration Monitoring Timer
After=network.target

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=service-orchestration-monitor.service

[Install]
WantedBy=timers.target 