# backups/version_2.4.0_beta_20250901_193155/templates/service-orchestrator.service.j2
[Unit]
Description=Service Orchestrator
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/service_orchestrator.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/service_orchestration/orchestrator.log
StandardError=append:/var/log/service_orchestration/orchestrator.log

[Install]
WantedBy=multi-user.target 