# backups/version_2.4.0_beta_20250901_193155/roles/media/templates/media-monitoring.service.j2
[Unit]
Description=Media Stack Monitoring Service
After=network.target

[Service]
Type=simple
User={{ username }}
ExecStart={{ config_dir }}/media-monitor.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target 