# templates/backup-orchestration.timer.j2
[Unit]
Description=Backup Orchestration Timer
Documentation=https://github.com/your-repo/ansible_homelab
Requires=backup-orchestration.service

[Timer]
# Run every hour at the top of the hour
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

# Accuracy settings
AccuracySec=1m

[Install]
WantedBy=timers.target 